# roles/ansible-lockdown.ubuntu24_cis/templates/etc/profile.d/50-umask.sh.j2
{{ file_managed_by_ansible }}

# Set umask with highest precedence

umask {{ ubtu24cis_bash_umask }}
