# roles/ansible-lockdown.ubuntu24_cis/templates/etc/systemd/system/dailyaidecheck.timer.j2
{{ file_managed_by_ansible }}

[Unit]
Description=Daily AIDE check

[Timer]
OnCalendar={{ ubtu24cis_aide_cron_day }}-{{ ubtu24cis_aide_cron_month }}-{{ ubtu24cis_aide_cron_weekday }} {{ ubtu24cis_aide_cron_hour }}:{{ ubtu24cis_aide_cron_minute }}:00
Persistent=true

[Install]
WantedBy=timers.target
